# buildroot/gadget/chippro/overlay/etc/inittab
# /etc/inittab
# Format for each entry: <id>:<runlevels>:<action>:<process>
#  id        == tty to run on, or empty for /dev/console
#  runlevels == ignored
#  action    == one of sysinit, respawn, askfirst, wait, and once
#  process   == program to run

null::sysinit:/etc/init.d/rcS &

ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100

::ctrlaltdel:/sbin/reboot

::shutdown:/etc/init.d/rcK
::shutdown:/bin/umount -a -r
